# examples/llm-proxy-gitagent/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
loglevel=info
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; ── claude-proxy ─────────────────────────────────────────────────────────────
; Starts FIRST (priority 5) — OAuth proxy to Anthropic API on :8319
; Reads token from /agentbox/.claude/.credentials.json (written by entrypoint)
[program:claude-proxy]
command=python3 /app/claude_proxy.py 8319
directory=/app
environment=
    HOME="/agentbox",
    PYTHONPATH="/app"
autostart=true
autorestart=true
startretries=10
startsecs=2
stopwaitsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=5

; ── ollama-facade ─────────────────────────────────────────────────────────────
; Starts SECOND (priority 8) — Ollama-compatible frontend on :11434
; Proxies to claude-proxy at 127.0.0.1:8319 via proxy_core.py
[program:ollama-facade]
command=python3 /app/ollama-facade/server.py
directory=/app
environment=
    HOME="/agentbox",
    PYTHONPATH="/app",
    LLM_PROXY_DB_PATH="/tmp/llm_proxy.db"
autostart=true
autorestart=true
startretries=10
startsecs=3
stopwaitsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=8

; ── OpenClaw Gateway ─────────────────────────────────────────────────────────
; Starts THIRD (priority 10) — AI agent runtime on :3000
[program:openclaw-gateway]
command=openclaw gateway run
directory=/agentbox
environment=
    HOME="/agentbox",
    NODE_ENV="production",
    OPENCLAW_CONFIG_PATH="/agentbox/.openclaw/openclaw.json",
    OLLAMA_API_KEY="ollama-local"
autostart=true
autorestart=true
startretries=5
startsecs=5
stopwaitsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=10

; ── backup-cron ──────────────────────────────────────────────────────────────
; Starts LAST (priority 15) — runs backup.sh every BACKUP_INTERVAL_HOURS hours
[program:backup-cron]
command=sh -c 'while true; do sleep $((${BACKUP_INTERVAL_HOURS:-6} * 3600)); /agentbox/bin/backup.sh; done'
directory=/agentbox
environment=
    HOME="/agentbox",
    BACKUP_INTERVAL_HOURS="%(ENV_BACKUP_INTERVAL_HOURS)s",
    BACKUP_KEEP_COUNT="%(ENV_BACKUP_KEEP_COUNT)s"
autostart=true
autorestart=true
startretries=3
startsecs=1
stopwaitsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=15
