# examples/ollama-llm-proxy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/agentbox/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/tmp/supervisord.pid
user=agentbox

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; ── 1. claude_proxy ──────────────────────────────────────────────────────────
; Starts first — OAuth proxy to Anthropic API on :8319
[program:claude-proxy]
command=python3 /app/llm-proxy/claude_proxy.py 8319
directory=/app/llm-proxy
environment=HOME="/agentbox"
autostart=true
autorestart=true
startretries=10
startsecs=3
stopwaitsecs=5
stdout_logfile=/agentbox/logs/claude-proxy.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_logfile=/agentbox/logs/claude-proxy.err
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=2
priority=5

; ── 2. ollama-facade ─────────────────────────────────────────────────────────
; Starts second — Ollama-compatible frontend on :11434, proxies to claude_proxy
[program:ollama-facade]
command=python3 /app/llm-proxy/ollama-facade/server.py
directory=/app/llm-proxy
environment=HOME="/agentbox",PYTHONPATH="/app/llm-proxy"
autostart=true
autorestart=true
startretries=10
startsecs=3
stopwaitsecs=5
stdout_logfile=/agentbox/logs/ollama-facade.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_logfile=/agentbox/logs/ollama-facade.err
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=2
priority=8

; ── 3. openclaw-gateway ──────────────────────────────────────────────────────
; Starts last — AI agent runtime on :3000, uses ollama-facade for LLM calls
[program:openclaw-gateway]
command=openclaw gateway run
directory=/agentbox
environment=HOME="/agentbox",NODE_ENV="production",OPENCLAW_CONFIG_PATH="/agentbox/.openclaw/openclaw.json"
autostart=true
autorestart=true
startretries=5
startsecs=5
stopwaitsecs=10
stdout_logfile=/agentbox/logs/openclaw.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/agentbox/logs/openclaw.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
priority=10
